# puckpulse-controller/packaging/puckpulse-update.service
[Unit]
Description=Check for PuckPulse Controller Updates
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/puckpulse-update
# Ensure we don't conflict with other apt processes
ExecStartPre=/usr/bin/sleep 10
